# rootfs/etc/systemd/system/libvirtd.service.d/override.conf
[Service]
ExecStartPre=/warewulf/sbin/execstartpre.%N
# ExecStartPost=/warewulf/sbin/execstartpost.%N.sh
